# deploy/kindshot-dashboard.service
[Unit]
Description=Kindshot Trading Dashboard (Streamlit)
After=network.target kindshot.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/kindshot
EnvironmentFile=/opt/kindshot/.env
ExecStart=/opt/kindshot/.venv/bin/streamlit run dashboard/app.py \
    --server.port 8501 \
    --server.address 0.0.0.0 \
    --server.headless true \
    --browser.gatherUsageStats false
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
